# service/supervisord/etc/conf.d/bench.conf
[program:bench]
command=./backend/target/debug/cryptocurrency run -c ./etc/validators/%(process_num)d.toml --leveldb-path ./vardb/%(process_num)d --port 85%(process_num)02d
process_name=%(program_name)s_%(process_num)02d
numprocs=9
autostart=false
autorestart=false
