# service/supervisord/etc/conf.d/frontend.conf
[program:frontend]
command=npm start -- --port=8280 --api-root=http://127.0.0.1:8200
process_name=%(program_name)s
numprocs=1
autostart=false
autorestart=false
