# service/supervisord/etc/conf.d/frontend.conf
[program:frontend]
command=node frontend/app.js --config=./etc/config.json
process_name=%(program_name)s
numprocs=1
autostart=false
autorestart=false
