# testnet/supervisord/etc/conf.d/drm_test.conf
[group:drm_test]
programs=drm_test

[program:drm_test]
command=digital_rights run -c drm_test/validators/0.toml --leveldb-path drm_test/db/0 --port 8600
autostart=false
autorestart=true