# testnet/supervisord/etc/conf.d/drm_backend.conf
[group:drm_backend]
programs=drm_backend

[program:drm_backend]
command=digital_rights --config drm_backend.conf run 0 --leveldb-path db/drm_backend_0 --port 8000
autostart=false
autorestart=true