# testnet/supervisord/etc/conf.d/drm_backend.conf
[group:drm_backend]
programs=drm_backend

[program:drm_backend]
command=digital_rights run -c drm_backend/validators/0.toml --leveldb-path drm_test/db/0 --port 8000
autostart=false
autorestart=true