# testnet/supervisord/etc/conf.d/drm_test.conf
[group:drm_test]
programs=drm_test

[program:drm_test]
command=digital_rights --config drm_test.conf run 0 --leveldb-path db/drm_test_0 --port 8600
autostart=false
autorestart=true